# lola-service/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=%(ENV_LOG_LEVEL)s

[program:libreoffice]
command=soffice --headless --norestore --nologo
    --env:UserInstallation=file:///tmp/lo_user_profile
    --accept="socket,host=localhost,port=%(ENV_LO_PORT)s;urp;StarOffice.ServiceManager"
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=100

[program:uvicorn]
command=uvicorn app.main:app --host 0.0.0.0 --port %(ENV_API_PORT)s --log-level %(ENV_LOG_LEVEL)s
directory=/app
autostart=false
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=200
